# supervisord3.conf
[supervisord]
nodaemon=true

[program:node3]
command=/opt/raiden/bin/raiden --config-file /opt/raiden/config/node3.toml --no-sync-check --user-deposit-contract-address %(ENV_USER_DEPOSIT_ADDRESS)s
