# supervisord4.conf
[supervisord]
nodaemon=true

[program:node4]
command=/opt/raiden/bin/raiden --config-file /opt/raiden/config/node4.toml --no-sync-check --user-deposit-contract-address %(ENV_USER_DEPOSIT_ADDRESS)s
